# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:wordpress]
command=/entrypoint.sh php-fpm -D

[program:test]
command=/run.sh
